# raspberrypi_exporter.timer
[Unit]
Description=Timer for the Raspberry Pi Prometheus exporter

[Timer]
OnBootSec=0min
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
